# systemd/zeroclaw@.service
[Unit]
Description=ZeroClaw daemon for user %i
Documentation=https://github.com/bighorse/zeroclaw
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# %h expands to the user's home directory (e.g. /home/claw-001).
WorkingDirectory=%h

# Per-instance environment file. Put ZEROCLAW_API_KEY etc. here so secrets
# never sit in config.toml. ClawOps writes /etc/clawops/zeroclaw-%i.env at
# provision time. The leading `-` keeps systemd happy when the file is absent.
EnvironmentFile=-/etc/clawops/zeroclaw-%i.env
EnvironmentFile=-/etc/clawops/zeroclaw.env

Environment=ZEROCLAW_CONFIG_DIR=%h/.zeroclaw
Environment=ZEROCLAW_WORKSPACE=%h/.zeroclaw/workspace
Environment=RUST_LOG=info,zeroclaw=info

ExecStart=/usr/local/bin/zeroclaw daemon --config-dir %h/.zeroclaw

Restart=on-failure
RestartSec=5

# Resource caps (work at user level)
MemoryMax=512M
TasksMax=64
LimitNOFILE=4096

# user-level supported hardening. Stronger directives like PrivateTmp,
# ProtectKernel*, ProtectControlGroups, RestrictAddressFamilies require
# capabilities the user systemd manager doesn't have and produce
# 218/CAPABILITIES. Per-user uid + 0700 home is the primary isolation.
NoNewPrivileges=true

[Install]
WantedBy=default.target
